# deploy/evo-server.service
[Unit]
Description=Evo-Server - Personal Programming Evolution Platform
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/evo-server
ExecStart=/usr/local/bin/python3 -m uvicorn evo_server.main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5
EnvironmentFile=/opt/evo-server/.env

[Install]
WantedBy=multi-user.target
